# nodejs/s6-overlay/etc/services.d/crond/run
#!/usr/bin/with-contenv bash

if [ -f /root/crontab.list ]; then
  crontab /root/crontab.list
fi

exec s6-setuidgid root crond -f
